# cron/update.cron
* * * * * curl lvs-web/start.php > /app/cron-out.txt
