# php/update.cron
* * * * * php /app/start.php > /app/cron-out.txt
